# st2/supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[inet_http_server]
port = 127.0.0.1:9001

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
nodaemon=true

[supervisorctl]
serverurl=http://127.0.0.1:9001
; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf

[program:st2actionrunner]
command=/opt/stackstorm/st2/bin/st2actionrunner --config-file /etc/st2/st2.conf
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2actionrunner.log
stderr_logfile = /var/log/st2actionrunner.log
user=root
numprocs=10
process_name=st2actionrunner%(process_num)s
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:st2api]
command=/opt/stackstorm/st2/bin/gunicorn_pecan /opt/stackstorm/st2/lib/python2.7/site-packages/st2api/gunicorn_config.py -k eventlet -b 127.0.0.1:9101 --workers 1 --threads 1 --graceful-timeout 10 --timeout 30
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2api.log
stderr_logfile = /var/log/st2api.log
user=root
numprocs=1
autorestart=true
startsecs=5
stopwaitsecs = 60
killasgroup=true

[program:st2auth]
command=/opt/stackstorm/st2/bin/gunicorn_pecan /opt/stackstorm/st2/lib/python2.7/site-packages/st2auth/gunicorn_config.py -k eventlet -b 127.0.0.1:9100 --workers 1 --threads 1 --graceful-timeout 10 --timeout 30
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2auth.log
stderr_logfile = /var/log/st2auth.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:st2garbagecollector]
command=/opt/stackstorm/st2/bin/st2garbagecollector --config-file /etc/st2/st2.conf
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2garbagecollector.log
stderr_logfile = /var/log/st2garbagecollector.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:st2notifier]
command=/opt/stackstorm/st2/bin/st2notifier --config-file /etc/st2/st2.conf
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2notifier.log
stderr_logfile = /var/log/st2notifier.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:st2resultstracker]
command=/opt/stackstorm/st2/bin/st2resultstracker --config-file /etc/st2/st2.conf
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2resultstracker.log
stderr_logfile = /var/log/st2resultstracker.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:st2sensorcontainer]
command=/opt/stackstorm/st2/bin/st2sensorcontainer --config-file /etc/st2/st2.conf
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2sensorcontainer.log
stderr_logfile = /var/log/st2sensorcontainer.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:st2rulesengine]
command=/opt/stackstorm/st2/bin/st2rulesengine --config-file /etc/st2/st2.conf
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2rulesengine.log
stderr_logfile = /var/log/st2rulesengine.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:st2stream]
command=/opt/stackstorm/st2/bin/gunicorn_pecan /opt/stackstorm/st2/lib/python2.7/site-packages/st2stream/gunicorn_config.py -k eventlet -b 127.0.0.1:9102 --workers 1 --threads 10 --graceful-timeout 10 --timeout 30
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/st2stream.log
stderr_logfile = /var/log/st2stream.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
redirect_stderr = true
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/nginx.log
stderr_logfile = /var/log/nginx.log
user=root
numprocs=1
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true

